# sway.target
[Unit]
Description=Sway Startup Target
Requires=default.target
After=default.target
